# templates/docker-image-prune.timer.j2
# Managed by Ansible — do not edit by hand.

[Unit]
Description=Weekly prune of unused Docker images

[Timer]
OnCalendar={{ docker_image_prune_oncalendar }}
# run on next boot if the host was off at the scheduled time
Persistent=true

[Install]
WantedBy=timers.target
